# 7.11-ubuntu-bionic/config/etc/supervisor/conf.d/redis.conf
[program:redis]
user=redis
command=/usr/bin/redis-server /etc/redis.conf
stopsignal=TERM
redirect_stderr=true
autorestart=true

