# 7.9/config/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
stopsignal=TERM
redirect_stderr=true
autorestart=true

